# deploy/archivist.service
[Unit]
Description=The Archivist — Personal Reading Tracker
After=network.target

[Service]
Type=simple
User=digitalhub
WorkingDirectory=/home/digitalhub/work/Personal-Reading-Tracker
EnvironmentFile=/home/digitalhub/work/Personal-Reading-Tracker/.env
ExecStart=/usr/bin/node /home/digitalhub/work/Personal-Reading-Tracker/node_modules/.bin/tsx server.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=archivist

[Install]
WantedBy=multi-user.target
